# ops/nginx/templates/etc-supervisor-conf.d-nginx.ini
[program:nginx]
command=/usr/sbin/nginx
user=root
autostart=true
autorestart=false
redirect_stderr=true
