# ops/acontrol-report/templates/etc-supervisor-conf.d-acontrol-report.ini
[program:report]
process_name=%(program_name)s
command=java -jar report-api/build/libs/report-api-0.1.0.jar
directory=/application-host
stdout_logfile=/var/log/reporting-api.log
user=root
autostart=false
autorestart=false
redirect_stderr=true
