# ops/spark-worker/templates/etc-supervisor-conf.d-spark-worker.ini
[program:spark-worker]
command=/usr/lib/spark-1.6.0-bin-hadoop2.6/sbin/start-slave.sh 127.0.0.1:7077
autostart=false
autorestart=false
user=root

