# ops/spark-worker/templates/supervisor-spark-worker.ini.j2
[program:spark-worker]
command={{ spark_lib_dir }}/spark-executor.sh
stdout_logfile={{ spark_log_dir }}/spark-worker.log
autostart=true
autorestart=false
user=root
killasgroup=true
stopasgroup=true
