# ops/acontrol/templates/etc-supervisor-conf.d-acontrol.ini
[program:import]
process_name=%(program_name)s
command=java -Xms128m -Xmx512m -ea -jar data-import/build/libs/data-import-0.1.0.jar
environment=SPARK_LOCAL_IP=10.1.0.100
directory=/application-host
stdout_logfile=/var/log/data-import.log
user=root
autostart=false
autorestart=false
redirect_stderr=true

[program:analytic]
process_name=%(program_name)s
command=java -Xms128m -Xmx512m -ea -jar data-analytic/build/libs/data-analytic-0.1.0.jar
directory=/application-host
stdout_logfile=/var/log/data-analytic.log
user=root
autostart=false
autorestart=false
redirect_stderr=true

[group:data]
programs=import,analytic
priority=10
