# ops/acontrol/templates/etc-supervisor-conf.d-acontrol.ini
[program:api-1]
process_name=%(program_name)s
command=java -jar /application-host/report-api/build/libs/report-api-0.1.0.jar --server.port=9901
stdout_logfile=/var/log/reporting-api-1.log
user=root
autostart=true
autorestart=false
redirect_stderr=true

[program:api-2]
process_name=%(program_name)s
command=java -jar /application-host/report-api/build/libs/report-api-0.1.0.jar --server.port=9902
stdout_logfile=/var/log/reporting-api-2.log
user=root
autostart=false
autorestart=false
redirect_stderr=true

[program:import]
process_name=%(program_name)s
command=java -jar /application-host/data-import/build/libs/data-import-0.1.0.jar
stdout_logfile=/var/log/data-import.log
user=root
autostart=false
autorestart=false
redirect_stderr=true

[program:analytic]
process_name=%(program_name)s
command=java -jar /application-host/data-analytic/build/libs/data-analytic-0.1.0.jar
stdout_logfile=/var/log/data-analytic.log
user=root
autostart=false
autorestart=false
redirect_stderr=true

[group:report]
programs=api-1,api-2
priority=9


[group:data]
programs=import,analytic
priority=10
